# salt/website/blog/files/blog-ipfs.path
[Unit]
Description=Publish blog.bricewge.fr on IPFS when it content changes

[Path]
PathChanged=/srv/nginx/blog.bricewge.fr

[Install]
WantedBy=multi-user.target