# filesystem/files/refind-update.service
[Unit]
Description=rEFInd bootloader update

[Service]
Type=oneshot
ExecStart=/usr/bin/refind-install