# salt/torrent/files/transmission-blocklist.timer
[Unit]
Description = "Daily update of Transmission's blocklist"

[Timer]
OnCalendar = Daily
Persistent = True

[Install]
WantedBy = timers.target
