# salt/filesystem/files/refind-update.path
[Unit]
Description=Path monitor for rEFInd updates

[Path]
PathChanged=/usr/share/refind

[Install]
WantedBy=multi-user.target