# salt/filesystem/btrfs-sxbackup/files/backup-attic.service
[Unit]
Description=Backup /srv/attic subvolume
After=local-fs.target

[Service]
ExecStart=/usr/local/bin/backup-attic.sh