# salt/http/nginx/nginx-reload.service
[Unit]
Description="Reload nginx configuration"

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload nginx.service