# salt/http/nginx/nginx-reload.path
[Unit]
Description=Reload nginx configuration when sites-available changes

[Path]
PathChanged=/etc/nginx/sites-available/blog.bricewge.fr
Unit=nginx-reload.path

[Install]
WantedBy=multi-user.target