# backup/files/borg.timer
[Unit]
Description = "Attic, daily backup"

[Timer]
OnCalendar = Daily
Persistent = True

[Install]
WantedBy = timers.target

