# torrent/files/flexget.timer
[Unit]
Description = "Flexget hourly"

[Timer]
OnCalendar = Hourly
Persistent = True

[Install]
WantedBy = timers.target

