# vm/files/cgmanager.service
[Unit]
Description=Control Group manager

[Service]
ExecStart=/usr/bin/cgmanager

[Install]
WantedBy=sysinit.target
